# systemd/neuforge-callcenter.service
[Unit]
Description=NeuForge Call Center IVR — BANT Flow on port 8002 (Flagship Commando)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/titan
ExecStart=/opt/titan/venv/bin/python /opt/titan/Titan_NeuForgeCallCenter.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/titan/neuforge_callcenter.log
StandardError=append:/var/log/titan/neuforge_callcenter.err
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/titan/.env
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
