# systemd/neuforge-sms.service
[Unit]
Description=NeuForge SMS Commando — 6-Touch Drip via Telnyx (Flagship Commando)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/titan
ExecStart=/opt/titan/venv/bin/python /opt/titan/Titan_NeuForgeSMS.py
Restart=always
RestartSec=30
StandardOutput=append:/var/log/titan/neuforge_sms.log
StandardError=append:/var/log/titan/neuforge_sms.err
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/titan/.env
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
